# supervisord.ini
[supervisord]
pidfile = /var/run/canalyzer.pid
logfile = /var/log/canalyzer.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
nodaemon = false
minfds = 1024
minprocs = 200

[unix_http_server]
file = /var/run/canalyzer_supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl = unix:///var/run/canalyzer_supervisord.sock


[program:uwsgi]
directory = /var/www/canalyzer.latest
command = venv/bin/uwsgi --ini /var/www/canalyzer.latest/uwsgi.ini
autostart = true
autorestart = true
stdout_logfile = /var/log/canalyzer_uwsgi.log
redirect_stderr = true
stopsignal = INT